# packer/webset.conf
[program:webset]
command=/home/webset/webset/run
autostart=true
autorestart=true
stderr_logfile=/var/log/webset.err.log
stdout_logfile=/var/log/webset.out.log
